# deploy/systemd/express-trade-kit-api.service
[Unit]
Description=Express Trade Kit API
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/express-trade-kit/current/server
EnvironmentFile=/var/www/express-trade-kit/current/server/.env
ExecStart=/usr/bin/node src/index.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

